# cfddns.service
[Unit]
Description=Cloudflare DDNS
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cfddns
ExecStart=/usr/bin/python3 /opt/cfddns/cfddns.py
Restart=always

[Install]
WantedBy=multi-user.target

